# templates/etc/profile.d/composer.sh.j2
# {{ ansible_managed }}

if [ -d "/root/.composer/vendor/bin" ]
    export PATH="$PATH:/root/.composer/vendor/bin"
fi
